# nupkg/payload/ProgramData/ManagedReports/queries-optimized.json
{
  "system_info": "SELECT hostname, cpu_brand, cpu_physical_cores, cpu_logical_cores, physical_memory FROM system_info;",
  
  "os_version": "SELECT name, version, build, platform, arch FROM os_version;",
  
  "hardware_info": "SELECT vendor, model, version, serial, hardware_vendor, hardware_model, hardware_version, hardware_serial FROM system_info;",
  
  "memory_info": "SELECT device, device_type, size, form_factor, configured_clock_speed FROM memory_devices WHERE size > 0;",
  
  "disk_info": "SELECT device, model, size, type FROM disk_info WHERE type IN ('SSD', 'HDD') AND size > 0;",
  
  "network_interfaces": "SELECT interface, mac, ip, mask, type, mtu FROM interface_details WHERE interface NOT LIKE 'lo%' AND interface NOT LIKE 'Teredo%' AND interface NOT LIKE '%Virtual%' AND mac != '';",
  
  "installed_programs": "SELECT name, version, publisher, install_date FROM programs WHERE name IS NOT NULL AND publisher IS NOT NULL AND name NOT LIKE 'Microsoft Visual C++%' AND name NOT LIKE 'KB%' ORDER BY install_date DESC LIMIT 200;",
  
  "services_critical": "SELECT name, display_name, status, start_type, path FROM services WHERE status = 'RUNNING' AND (name LIKE '%defender%' OR name LIKE '%firewall%' OR name LIKE '%antivirus%' OR name LIKE '%security%' OR name LIKE '%winrm%' OR name LIKE '%rpc%' OR name LIKE '%dns%' OR name LIKE '%dhcp%');",
  
  "startup_items": "SELECT name, path, source, status, username FROM startup_items WHERE status = 'enabled';",
  
  "windows_crashes": "SELECT datetime, module, path, pid FROM windows_crashes WHERE datetime > datetime('now', '-7 days') LIMIT 10;",
  
  "patches_recent": "SELECT csname, hotfix_id, description, installed_on FROM patches WHERE installed_on > date('now', '-30 days') ORDER BY installed_on DESC LIMIT 20;",
  
  "bitlocker_info": "SELECT device_id, drive_letter, conversion_status, protection_status, lock_status, encryption_method FROM bitlocker_info;",
  
  "windows_security_center": "SELECT domain, antivirus, antispyware, firewall, autoupdate, uac FROM windows_security_center;",
  
  "tpm_info": "SELECT activated, enabled, owned, manufacturer_version, tpm_version_major, tpm_version_minor FROM tpm_info;",
  
  "firewall_rules_enabled": "SELECT name, direction, action, protocol, local_port, remote_port FROM windows_firewall_rules WHERE enabled = 1 AND action = 'Allow' AND (local_port != '' OR remote_port != '') LIMIT 50;",
  
  "certificates_important": "SELECT common_name, subject, issuer, not_valid_before, not_valid_after, key_algorithm FROM certificates WHERE (path = 'CurrentUser/My' OR path = 'LocalMachine/My') AND not_valid_after > datetime('now') LIMIT 25;",
  
  "logged_in_users": "SELECT user, host, time, pid FROM logged_in_users WHERE user != 'SYSTEM' AND user != 'LOCAL SERVICE' AND user != 'NETWORK SERVICE';",
  
  "user_groups": "SELECT uid, gid, username, description FROM user_groups WHERE username IS NOT NULL AND gid < 1000;",
  
  "processes_notable": "SELECT pid, name, path, cmdline, parent, threads, disk_bytes_read, disk_bytes_written FROM processes WHERE (name LIKE '%explorer%' OR name LIKE '%chrome%' OR name LIKE '%firefox%' OR name LIKE '%edge%' OR name LIKE '%defender%' OR name LIKE '%antivirus%' OR name LIKE '%malware%' OR cmdline LIKE '%powershell%' OR cmdline LIKE '%cmd%') AND name NOT LIKE 'svchost%' LIMIT 50;",
  
  "listening_ports": "SELECT pid, port, protocol, family, address, path FROM listening_ports WHERE port != 0 AND address NOT LIKE '127.0.0.1' AND address NOT LIKE '::1' LIMIT 30;",
  
  "scheduled_tasks_custom": "SELECT name, action, path, enabled, state, last_run_time, next_run_time FROM scheduled_tasks WHERE enabled = 1 AND path NOT LIKE '\\Microsoft\\%' LIMIT 25;",
  
  "registry_security": "SELECT key, name, type, data FROM registry WHERE key LIKE 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System%' AND name IN ('EnableLUA', 'ConsentPromptBehaviorAdmin', 'PromptOnSecureDesktop');",
  
  "windows_defender": "SELECT anti_virus, anti_spyware, real_time_protection, on_access_protection, network_protection, behavior_monitor FROM windows_security_products WHERE type = 'Antivirus' LIMIT 5;",
  
  "usb_devices": "SELECT vendor, model, serial, removable, size FROM usb_devices WHERE removable = 1 AND size > 0;",
  
  "browser_plugins_security": "SELECT name, version, description, path FROM browser_plugins WHERE (browser = 'chrome' OR browser = 'firefox' OR browser = 'edge') AND (name LIKE '%security%' OR name LIKE '%antivirus%' OR name LIKE '%ad%block%' OR name LIKE '%password%') LIMIT 15;",
  
  "event_logs_system": "SELECT datetime, source, provider_name, eventid, level FROM windows_events WHERE channel = 'System' AND level <= 3 AND datetime > datetime('now', '-24 hours') AND eventid IN (1074, 6005, 6006, 6008, 6009, 6013) ORDER BY datetime DESC LIMIT 10;",
  
  "event_logs_security": "SELECT datetime, source, provider_name, eventid, level FROM windows_events WHERE channel = 'Security' AND eventid IN (4624, 4625, 4648, 4672, 4720, 4732, 4756) AND datetime > datetime('now', '-24 hours') ORDER BY datetime DESC LIMIT 15;",
  
  "wmi_cli_event_consumers": "SELECT name, command_line_template FROM wmi_cli_event_consumers;",
  
  "autoexec": "SELECT name, path, source FROM autoexec WHERE path IS NOT NULL;",
  
  "file_events_startup": "SELECT action, target_path, md5, atime, mtime FROM file_events WHERE action = 'CREATED' AND target_path LIKE '%\\Startup\\%' AND target_path LIKE '%.exe' LIMIT 5;"
}
