# src/osquery-queries.json
{
  "system_info": "SELECT hostname, cpu_brand, cpu_physical_cores, cpu_logical_cores, physical_memory, hardware_vendor, hardware_model, hardware_serial, uuid FROM system_info;",
  
  "os_version": "SELECT name, version, build, platform, arch, install_date FROM os_version;",
  
  "os_details": "SELECT install_date, codename FROM os_version;",
  
  "windows_experience_pack": "SELECT key, name, type, data FROM registry WHERE key = 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion' AND name IN ('UBR', 'DisplayVersion', 'ReleaseId', 'CurrentBuild', 'InstallDate') AND data IS NOT NULL;",
  
  "hardware_info": "SELECT hardware_vendor, hardware_model, hardware_version, hardware_serial, uuid FROM system_info;",
  
  "cpu_info": "SELECT device_id, name, manufacturer, family, model, stepping, processor_type, characteristics, address_width, current_clock_speed, max_clock_speed, cores, logical_processors FROM cpu_info;",
  
  "processor_registry": "SELECT key, name, type, data FROM registry WHERE key = 'HKEY_LOCAL_MACHINE\\HARDWARE\\DESCRIPTION\\System\\CentralProcessor\\0' AND name IN ('ProcessorNameString', 'Identifier', 'VendorIdentifier') AND data IS NOT NULL;",
  
  "processor_all_registry": "SELECT key, name, type, data FROM registry WHERE key LIKE 'HKEY_LOCAL_MACHINE\\HARDWARE\\DESCRIPTION\\System\\CentralProcessor\\%' AND name = 'ProcessorNameString' AND data IS NOT NULL;",
  
  "wmi_processor": "SELECT name, manufacturer, description, max_clock_speed, current_clock_speed, number_of_cores, number_of_logical_processors, architecture, family, model, stepping FROM wmi_cpu_info;",
  
  "system_board": "SELECT serial_number, manufacturer, product, version FROM system_info;",
  
  "memory_info": "SELECT device_locator, memory_type, size, form_factor, configured_clock_speed, manufacturer, part_number FROM memory_devices WHERE size > 0;",
  
  "disk_info": "SELECT id, name, disk_index, disk_size, type, manufacturer, model, serial_number, description FROM disk_info WHERE disk_size > 0;",
  
  "video_info": "SELECT name, driver_version, driver_date, adapter_ram, current_horizontal_resolution, current_vertical_resolution, adapter_compatibility FROM video_info;",
  
  "network_interfaces": "SELECT interface, mac, type, mtu, friendly_name, description, manufacturer, connection_id, dhcp_enabled, dhcp_server FROM interface_details WHERE interface NOT LIKE 'lo%' AND interface NOT LIKE 'Loopback%';",
  
  "network_addresses": "SELECT interface, address, mask, broadcast, type FROM interface_addresses WHERE address IS NOT NULL AND address NOT LIKE '127.%' AND address NOT LIKE '::1%';",
  
  "wifi_networks": "SELECT ssid, bssid, network_name, security_type, signal_strength, authentication FROM wifi_networks WHERE ssid IS NOT NULL;",
  
  "route_table": "SELECT destination, netmask, gateway, source, flags, interface, mtu, metric, type FROM routes WHERE destination IS NOT NULL;",
  
  "installed_programs": "SELECT name, version, publisher, install_date FROM programs WHERE name IS NOT NULL AND name != '';",
  
  "running_services": "SELECT name, display_name, status, start_type, path FROM services WHERE status = 'RUNNING' AND name IS NOT NULL;",
  
  "startup_items": "SELECT name, path, source, status, username FROM startup_items WHERE name IS NOT NULL;",
  
  "recent_patches": "SELECT csname, hotfix_id, description, installed_on FROM patches WHERE installed_on IS NOT NULL ORDER BY installed_on DESC LIMIT 20;",
  
  "bitlocker_info": "SELECT device_id, drive_letter, conversion_status, protection_status, lock_status, encryption_method FROM bitlocker_info WHERE device_id IS NOT NULL;",
  
  "security_center": "SELECT antivirus, antispyware, firewall, autoupdate FROM windows_security_center;",
  
  "windows_defender": "SELECT key, name, type, data FROM registry WHERE key = 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows Defender' AND name IN ('InstallLocation', 'ProductVersion') AND data IS NOT NULL;",
  
  "defender_exclusions": "SELECT key, name, type, data FROM registry WHERE key LIKE 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows Defender\\Exclusions\\%' AND data IS NOT NULL;",
  
  "windows_firewall": "SELECT key, name, type, data FROM registry WHERE key LIKE 'HKEY_LOCAL_MACHINE\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\%Profile' AND name IN ('EnableFirewall', 'DefaultInboundAction', 'DefaultOutboundAction') AND data IS NOT NULL;",
  
  "tpm_info": "SELECT activated, enabled, owned, manufacturer_version, manufacturer_id, manufacturer_name, spec_version FROM tpm_info;",
  
  "mdm_enrollment": "SELECT enrollment_id, enrollment_type, enrollment_state, management_service_url, certificate_thumbprint FROM mdm_enrollment;",
  
  "enrollment_info": "SELECT key, name, type, data FROM registry WHERE key LIKE 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Enrollments\\%' AND name IN ('ProviderID', 'UPN', 'EnrollmentState', 'AADResourceID', 'EnrollmentType') AND data IS NOT NULL;",
  
  "critical_processes": "SELECT pid, name, path, cmdline, parent, threads FROM processes WHERE name IN ('lsass.exe', 'winlogon.exe', 'csrss.exe', 'smss.exe', 'explorer.exe') AND name IS NOT NULL;",
  
  "listening_ports": "SELECT pid, port, protocol, family, address, path FROM listening_ports WHERE port != 0 AND port IS NOT NULL;",
  
  "active_tasks": "SELECT name, action, path, enabled, state, last_run_time, next_run_time FROM scheduled_tasks WHERE enabled = 1 AND name IS NOT NULL LIMIT 30;",
  
  "security_policies": "SELECT key, name, type, data FROM registry WHERE key = 'HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System' AND name IN ('EnableLUA', 'ConsentPromptBehaviorAdmin', 'PromptOnSecureDesktop') AND data IS NOT NULL;",
  
  "logical_drives": "SELECT device_id, type, description, size, free_space, file_system FROM logical_drives WHERE device_id IS NOT NULL;",
  
  "logged_in_users": "SELECT user, host, time, pid FROM logged_in_users WHERE user IS NOT NULL;",
  
  "system_uptime": "SELECT total_seconds FROM uptime;",
  
  "battery_info": "SELECT serial_number, manufacturer, chemistry, designed_capacity, max_capacity, current_capacity, percent, amperage, voltage, cycle_count, condition FROM battery WHERE serial_number IS NOT NULL;"
}
